# config/jupiter/roles/hosts/templates/etc/systemd/system/sport_sync.timer
[Unit]
Description=Run sport sync every 12 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h 

[Install]
WantedBy=timers.target
